# tests/07-mount-point/run.expected/units/systemd-zram-setup@zram12.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=var-folded.mount
