# tests/01-basic/run.expected/units/swap-create@.service
# Automatically generated by zram-generator

[Unit]
Description=Create swap on /dev/%i
Wants=systemd-modules-load.service
After=systemd-modules-load.service
After=dev-%i.device
DefaultDependencies=false

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-modprobe zram
ExecStart=/path/to/zram-generator --setup-device '%i'
