# tests/07-mount-point/run.expected/units/var-compressed.mount
# Automatically generated by "/home/zbyszek/src/zram-generator/target/debug/zram-generator"

[Unit]
Description=Compressed Storage on /dev/zram11
Documentation=man:zram-generator(8) man:zram-generator.conf(5)
Requires=systemd-zram-setup@zram11.service
After=systemd-zram-setup@zram11.service

[Mount]
What=/dev/zram11
Where="/var/compressed"
