# units/systemd-zram-setup@.service.d/binary-location.conf
# SPDX-License-Identifier: MIT
# This file is part of the zram-generator project
# https://github.com/systemd/zram-generator

[Service]
ExecStart=
ExecStart={generator} --setup-device '%i'
