# tests/07-mount-point/run.expected/units/systemd-zram-setup@zram15.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=-.mount
