# tests/09-zram-size/run.expected/units/systemd-zram-setup@zram0.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=dev-%i.swap
