# tests/10-example/run.expected/units/systemd-zram-setup@zram1.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=run-compressed\x2dmount\x2dpoint.mount
