# units/swap-create@.service.d/binary-location.conf
# SPDX-License-Identifier: MIT

[Service]
ExecStart=
ExecStart={generator} --setup-device '%i'
