# tests/07-mount-point/run.expected/units/systemd-zram-setup@zram13.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=foo-bar-baz.mount
