# tests/10-example/run.expected/units/run-compressed\x2dmount\x2dpoint.mount
# Automatically generated by zram-generator

[Unit]
Description=Compressed Storage on /dev/zram1
Documentation=man:zram-generator(8) man:zram-generator.conf(5)
Requires=systemd-zram-setup@zram1.service
After=systemd-zram-setup@zram1.service

[Mount]
What=/dev/zram1
Where=/run/compressed-mount-point
Options=discard
