# tests/10-example/run.expected/units/systemd-zram-setup@zram0.service.d/bindings.conf
# Automatically generated by zram-generator

[Unit]
BindsTo=dev-%i.swap dev-zvol-tarta\x2dzoot-swap\x2dwriteback.device
After=dev-zvol-tarta\x2dzoot-swap\x2dwriteback.device
